# rootfs/etc/services.d/sdrplay/run
#!/bin/bash

/usr/local/bin/sdrplay_apiService
