# rootfs/etc/services.d/rtl_tcp/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

if [[ -n "$SOAPY" ]]; then
    python /scripts/soapy2tcp.py
else
    echo "SOAPY env var must be set"
    sleep 86400
fi
