# rootfs/etc/services.d/00-rtl_tcp/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

rtl_tcp -g $GAIN -s $RATE -d $SERIAL -f $FREQ
